# resources/services.d/nginx/finish
#!/bin/sh

echo stopping nginx...

exec s6-svscanctl -t /var/run/s6/services
