# resources/services.d/nginx/finish
#!/usr/bin/execlineb -S0

foreground { echo stopping nginx... }

s6-svscanctl -t /var/run/s6/services
