# resources/services.d/nginx/run
#!/usr/bin/execlineb -P

foreground { echo starting nginx... }

nginx -c /etc/nginx/nginx.conf -g "daemon off;"
