# resources/services.d/nginx/run
#!/bin/sh

echo starting nginx...

exec nginx -c /etc/nginx/nginx.conf -g "daemon off;"
